# deploy/claude-peers-broker.service
[Unit]
Description=Claude Peers Broker
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=jamditis
ExecStart=/home/peer/.bun/bin/bun /home/peer/projects/claude-peers-mcp/broker.ts
Restart=always
RestartSec=5
Environment=HOME=/home/peer

[Install]
WantedBy=multi-user.target
