# deploy/claude-peers-broker.service
[Unit]
Description=Claude Peers Broker
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=jamditis
ExecStart=/home/jamditis/.bun/bin/bun /home/jamditis/projects/claude-peers-mcp/broker.ts
Restart=always
RestartSec=5
Environment=HOME=/home/jamditis
# This broker is supervised; never let it self-exit when idle or it would restart-loop.
Environment=CLAUDE_PEERS_IDLE_EXIT_MS=0

[Install]
WantedBy=multi-user.target
